# solutions/board/rootfs-overlay/etc/udev/rules.d/99-qemu-usb-kbd.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0627", ATTR{idProduct}=="0001", RUN+="/sbin/modprobe solution_1"
